# deploy/apsales-whatsapp-bridge.service.d/crash-logger.conf
[Service]
ExecStopPost=/usr/bin/node /root/.openclaw/scripts/apsales-bridge-crash-logger.mjs
