# ops/launchd/ai.asiapower.apcoo-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.asiapower.apcoo-bot</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/longhui/Desktop/AsiaPower/.venv/bin/python</string>
        <string>/Users/longhui/Desktop/AsiaPower/integrations/telegram_coo_bot.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/longhui/Desktop/AsiaPower</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>StandardOutPath</key>
    <string>/Users/longhui/Desktop/AsiaPower/memory/apcoo-bot.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/longhui/Desktop/AsiaPower/memory/apcoo-bot.log</string>
</dict>
</plist>
