# deploy/launchd/com.asiapower.apbd.continuous.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.asiapower.apbd.continuous</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/longhui/Desktop/AsiaPower/.venv/bin/python3</string>
    <string>/Users/longhui/Desktop/AsiaPower/scripts/apbd-continuous.py</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    <key>APBD_RUNNER_INTERVAL_SECONDS</key>
    <string>2700</string>
  </dict>
  <key>WorkingDirectory</key>
  <string>/Users/longhui/Desktop/AsiaPower</string>
  <key>StandardOutPath</key>
  <string>/Users/longhui/Desktop/AsiaPower/reports/apbd-continuous.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/longhui/Desktop/AsiaPower/reports/apbd-continuous.log</string>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>60</integer>
  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
