# deploy/apbd-ca-leads-trickle.service
[Unit]
Description=APBD Canada auto-repair leads trickle (load-gated continuous)
After=network.target
Documentation=file:///root/.openclaw/workspace/AsiaPower/docs/agents/apbd/lead-discovery.md

[Service]
Type=simple
WorkingDirectory=/root/.openclaw/workspace/AsiaPower
EnvironmentFile=-/root/.openclaw/workspace/AsiaPower/.env
Environment=PYTHONUNBUFFERED=1
# Keep collecting while load is OK until free Places quota is exhausted (no 4h cron).
ExecStart=/root/.openclaw/workspace/AsiaPower/.venv/bin/python3 scripts/apbd_leads_ca_trickle.py --loop --limit 8 --max-cities 1 --load-max 1.8 --idle-sleep 20 --busy-sleep 60 --quota-sleep 3600
StandardOutput=journal
StandardError=journal
Restart=always
RestartSec=60
Nice=15
IOSchedulingClass=best-effort
IOSchedulingPriority=7
# Keep this worker small — must not crowd inventory-site / WhatsApp
MemoryMax=256M
MemoryHigh=192M
CPUQuota=35%
TasksMax=32

[Install]
WantedBy=multi-user.target
